# systemd/bcs-update.timer
[Unit]
Description=Build Check Statistics Daily Update Timer

[Timer]
OnCalendar=02:00:00
