# systemd/bcs-update.timer
[Unit]
Description=Build Check Statistics update timer

[Timer]
OnCalendar=02:00:00
